# etc/talkpile.service.plist
<!-- Edit below as needed and save to ~/Library/LaunchAgents/talkpile.service.plist -->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>talkpile.service</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/your-user-name/Library/Logs/talkpile.service.stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/your-user-name/Library/Logs/talkpile.service.stdout.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>FORCE_COLOR</key>
      <string>true</string>
      <key>PATH</key>
      <string><![CDATA[/Users/your-user-name/n/bin:/opt/homebrew/bin:/usr/local/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin]]></string>
      <key>OPENAI_API_KEY</key>
      <string><![CDATA[your-api-key]]></string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/your-user-name/repos/github.com/machellerogden/talkpile</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/your-user-name/n/bin/node</string>
      <string>bin/service.mjs</string>
      <string>--debug</string>
      <!--<string>--verbose</string>-->
    </array>
  </dict>
</plist>
